# .devcontainer/devcron/crontab
15 21 * * * /usr/local/bin/update_data.sh >> /workspace/.logs/update_data.log 2>&1
15 21 * * * /usr/local/bin/run_pipeline.sh >> /workspace/.logs/pipeline.log 2>&1
